# fedora_update_job.crontab
 0 17 * * * sudo dnf distro-sync -y >> fedora_update_jog_log.txt
