# update.crontab
* * * * * sudo pacmatic -Syu >> ~/.pacmatic_update.log
